# default/4MB/etc/rc.local
/etc/init.d/macchanger start
/etc/init.d/rt3050-event-tracker /dev/kmsg &
exit 0