# default/4MB/etc/hotplug.d/usb/30-usb-drive
#!/bin/bash

[[ $ACTION == "add" ]] && {
    echo "USB Drive plugged." > /tmp/usb_action
} || {
    echo "USB Drive unplugged." > /tmp/usb_action
}