# wmr-300/etc/hotplug.d/iface/75-jabberjaw
#!/bin/sh

[ "$ACTION" == "ifup" ] && [ "$INTERFACE" == "wan" ] && {
    logger -t DEBUG "New Ethernet connection detected. Execute payload...";
    sh /root/payload/payload.sh
}
